# control-plane/terraform/modules/site/rootfs/etc/systemd/system/reconcile-secrets.service
[Unit]
Description=Secrets Reconciliation Service
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/etc/default/reconcile-secrets
ExecStart=/usr/local/bin/reconcile-secrets.sh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=reconcile-secrets

[Install]
WantedBy=multi-user.target
