# control-plane/terraform/modules/site/rootfs/etc/systemd/system/watch-rollout-trigger.service
[Unit]
Description=Rollout Trigger Watcher
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=-/etc/default/reconcile-rollout
ExecStart=/usr/local/bin/watch-rollout-trigger.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=watch-rollout-trigger

[Install]
WantedBy=multi-user.target
